# etc/systemd/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin owbosh --noclear %I $TERM
Type=simple
